# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nextjs]
command=node /app/server.js
environment=PORT="3001",NODE_ENV="production",DATABASE_URL="file:/data/zenith.db"
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=5
startretries=10

[program:chat-service]
command=bun run /app/mini-services/chat-service/index.ts
directory=/app/mini-services/chat-service
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=3
startretries=10

[program:tcp-service]
command=bun run /app/mini-services/tcp-service/index.ts
directory=/app/mini-services/tcp-service
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=3
startretries=10

[program:udp-service]
command=bun run /app/mini-services/udp-service/index.ts
directory=/app/mini-services/udp-service
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=3
startretries=10

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=3
startretries=10